# error_record/20260207/logrotate/crontab
0 2 * * * sudo /usr/sbin/logrotate -v -s /opt/project/PAP/logrotate/logrotate.state /opt/project/PAP/logrotate/logrotate.conf > /opt/project/PAP/logrotate/logrotate.log 2>&1